# hsctl/systemd/hsctl-backup.service
[Unit]
Description=hsctl homeserver backup (restic)
After=docker.service
Wants=docker.service
# Don't even start the backup unless the external HDD is mounted (systemd pulls in the
# .mount unit for this path and waits for it). Keep in sync with REQUIRE_MOUNT/RESTIC_REPO
# in backup.conf. hsctl also enforces this in-process via the REQUIRE_MOUNT guard.
RequiresMountsFor=/mnt/restic

[Service]
Type=oneshot
# Root so restic can read Docker volume files under /var/lib/docker/volumes.
# Replace __DIR__ with the repo path.
WorkingDirectory=__DIR__
ExecStart=/usr/local/bin/hsctl backup run
User=root
